# dist/macos/bundle/Synergy.app/Contents/Resources/com.symless.synergy.v2.synergyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
        <string>com.symless.synergy.v2.synergyd</string>
    <key>LimitLoadToSessionType</key>
        <array>
            <string>System</string>
            <string>Aqua</string>
        </array>
    <key>KeepAlive</key>
    <true/>
        <key>ProgramArguments</key>
    <array>
            <string>/Applications/Synergy.app/Contents/MacOS/synergyd</string>
    </array>
</dict>
</plist>
